# bootc/services/mongodb/rootfs/usr/lib/systemd/system/mongod.service.d/override.conf
[Service]
ExecStartPre=/usr/libexec/testboot/gen-password.sh /var/lib/mongodb/.admin-pw 32
Environment="OPTIONS=-f /usr/share/mongodb/mongod.conf"
StateDirectory=mongodb
LogsDirectory=mongodb
