# bootc/services/mongodb/rootfs/usr/lib/systemd/system/mongod.service.d/override.conf
[Service]
StateDirectory=mongodb
LogsDirectory=mongodb
