# bootc/services/mongodb/rootfs/usr/lib/systemd/system/mongod.service.d/override.conf
[Unit]
After=mongodb-setup.service

[Service]
Environment="OPTIONS=-f /usr/share/mongodb/mongod.conf"
StateDirectory=mongodb
LogsDirectory=mongodb
# Relabel /var/lib/mongodb with the fcontext rule baked at build time.
# Runs before every start; safe because restorecon is idempotent.
# '+' prefix: run as root (upstream mongod.service sets User=mongod).
ExecStartPre=+/usr/sbin/restorecon -Rv /var/lib/mongodb
