# bootc/services/mongodb/rootfs/usr/lib/systemd/system/mongod.service.d/override.conf
[Unit]
After=mongodb-setup.service

[Service]
Environment="OPTIONS=-f /usr/share/mongodb/mongod.conf"
StateDirectory=mongodb
LogsDirectory=mongodb
