# bootc/apps/hello/rootfs/usr/lib/systemd/system/hello.service
[Unit]
Description=Hello World HTTP Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hello
Group=hello
StateDirectory=bootc-testboot/hello
LogsDirectory=bootc-testboot/hello
Environment=LOG_FILE=/var/log/bootc-testboot/hello/hello.log
# LOG_LEVEL=info LOG_FORMAT=text|json — see repos/hello
ExecStart=/usr/bin/hello
ExecStartPost=/usr/bin/env LOG_FILE=/var/log/bootc-testboot/hello/healthcheck.log /usr/libexec/testboot/healthcheck.sh http://127.0.0.1:8080/health 10
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
