# bootc/services/mongodb/rootfs/usr/lib/systemd/system/mongodb-setup.service
[Unit]
Description=MongoDB credential and TLS setup
Before=mongod.service
ConditionPathExists=!/var/lib/mongodb/.setup-done

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/testboot/mongodb-setup.sh

[Install]
WantedBy=multi-user.target
