# bootc/services/valkey/rootfs/usr/lib/systemd/system/valkey.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/valkey-server /usr/share/valkey/valkey.conf --daemonize no --supervised systemd
StateDirectory=valkey
LogsDirectory=valkey
