# bootc/services/rabbitmq/rootfs/usr/lib/systemd/system/rabbitmq-server.service.d/override.conf
[Service]
StateDirectory=rabbitmq
LogsDirectory=rabbitmq
