# bootc/services/redis/rootfs/usr/lib/systemd/system/redis.service.d/override.conf
[Service]
StateDirectory=redis
LogsDirectory=redis
