# bootc/services/redis/rootfs/usr/lib/systemd/system/redis.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/redis-server /usr/share/redis/redis.conf --daemonize no --supervised systemd
StateDirectory=redis
LogsDirectory=redis
